# eternal-daemon.service
[Unit]
Description=Eternal Daemon
[Service]
User=root
ExecStart=/usr/local/bin/eternal-daemon
Restart=on-failure
[Install]
WantedBy=multi-user.target